# 石论后台.app/Contents/Resources/cn.shilun.backend.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>cn.shilun.backend</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/shibicheng/shilun_standalone/.venv312/bin/python</string>
        <string>-m</string>
        <string>shilun</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/shibicheng/shilun_standalone</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ProcessType</key>
    <string>Background</string>
    <key>ThrottleInterval</key>
    <integer>5</integer>
    <key>StandardOutPath</key>
    <string>/Users/shibicheng/shilun_standalone/runtime/shilun-backend.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/shibicheng/shilun_standalone/runtime/shilun-backend.log</string>
</dict>
</plist>
